# src/connman.service.in
[Unit]
Description=Connection service
After=syslog.target

[Service]
EnvironmentFile=-/etc/sysconfig/connman
Type=dbus
BusName=net.connman
Restart=on-failure
ExecStart=@prefix@/sbin/connmand -n $OPTIONS
StandardOutput=null

[Install]
WantedBy=multi-user.target
